# metro_display/systemd/metro-display.service
[Unit]
Description=Metro Display Desktop
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/m4_marmaray_timeline
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONPATH=/home/pi/Desktop/m4_marmaray_timeline
ExecStart=/home/pi/Desktop/m4_marmaray_timeline/.venv/bin/python -m metro_display.desktop
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
